# simple_camio.service
[Unit]
Description=Simple CamIO - Interactive Tactile Map System
After=network.target sound.target

[Service]
Type=simple
User=user
Group=user
WorkingDirectory=/home/user/simple_camio

# Environment variables for headless operation
# DISPLAY required for Pyglet audio (even without actual display)
Environment="DISPLAY=:0"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse/"

# Run in headless mode (no display window)
# Using xvfb-run for virtual display support (recommended)
# Install with: sudo apt-get install xvfb
ExecStart=/usr/bin/xvfb-run -a /home/user/simple_camio/venv/bin/python /home/user/simple_camio/simple_camio.py --headless --input1 /home/user/simple_camio/models/CnapMap/CnapFirstFloor.json

# Alternative without xvfb (if DISPLAY and audio work):
# ExecStart=/home/user/simple_camio/venv/bin/python /home/user/simple_camio/simple_camio.py --headless --input1 /home/user/simple_camio/models/CnapMap/CnapFirstFloor.json

# Restart on failure
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=simple_camio

# Security settings (optional, adjust as needed)
# NoNewPrivileges=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target
